# workshop/supervisor/jupyter.conf
[program:jupyter]
process_name=jupyter
autorestart=true
command=/opt/conda/start-jupyter.sh
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
