# workshop/supervisor/jupyter.conf
[program:jupyter]
process_name=jupyter
command=jupyter lab --ip 0.0.0.0 --port 8888
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
