# workshop-images/base-environment/opt/eduk8s/etc/profile.d/01-kubernetes.sh
# Read in Kubernetes environment variables set from setup script.

if [ -f $HOME/.local/share/workshop/kubernetes-settings.sh ]; then
    . $HOME/.local/share/workshop/kubernetes-settings.sh
fi
