# workshop-images/base-environment/opt/eduk8s/etc/profile.d/02-webdav.sh
# Read in WebDAV environment variables set from setup script.

if [ -f $HOME/.local/share/workshop/webdav-settings.sh ]; then
    . $HOME/.local/share/workshop/webdav-settings.sh
fi
