# workshop-images/base-environment/opt/eduk8s/etc/supervisor/git.conf
[program:git]
process_name=git
command=/opt/eduk8s/sbin/start-git
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autostart=%(ENV_ENABLE_GIT_PROCESS)s
autorestart=true
