# workshop-images/base-environment/opt/eduk8s/etc/supervisor/httpd.conf
[program:httpd]
process_name=httpd
command=/opt/eduk8s/sbin/start-httpd
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autostart=%(ENV_ENABLE_HTTPD_PROCESS)s
autorestart=true
