# dashboard/opt/eduk8s/etc/supervisor/terminal.conf
[program:terminal]
command=/opt/eduk8s/sbin/start-butterfly
stopsignal=KILL
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autostart=%(ENV_ENABLE_TERMINAL)s
