# dashboard/opt/workshop/etc/supervisor/console.conf
[program:console]
process_name=console
command=/opt/workshop/bin/start-console.sh
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autostart=%(ENV_ENABLE_CONSOLE)s
