# dashboard/opt/workshop/etc/supervisor/gateway.conf
[program:gateway]
command=/opt/workshop/sbin/start-gateway
stopsignal=KILL
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
