# dashboard/etc/profile.d/eduk8s.sh
if [ -f /opt/eduk8s/.bash_profile ]; then
    . /opt/eduk8s/.bash_profile
fi
