# portal/opt/eduk8s/etc/supervisor/portal.conf
[program:portal]
process_name=portal
command=/opt/eduk8s/sbin/start-portal
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
